# gateway/systemd-files/smarthome-gateway.service.d/smarthome-ipv6.conf
[Unit]
After=ip6tables.service

[Service]
ExecStartPre=/sbin/ip6tables -w -A INPUT -p tcp --dport 8080 -j ACCEPT
ExecStopPost=/sbin/ip6tables -w -D INPUT -p tcp --dport 8080 -j ACCEPT
